# services/external_devices/run
#!/bin/sh
echo "*** starting external devices service ***"
exec 2>&1
exec /data/external-devices/external_devices.py
