# services/virtual_devices/run
#!/bin/sh
echo "*** starting MQTT switches service ***"
exec 2>&1
exec /data/venus-os_virtual-devices/virtual_devices.py
