# services/mqtt_switches/run
#!/bin/sh
echo "*** starting MQTT switches service ***"
exec 2>&1
exec /data/MQTT-switches/mqtt_switches.py
